# files/usr/lib/systemd/system/bootc-fetch-apply-updates.timer.d/updates.conf
[Timer]
# Clear previous timers
OnBootSec= OnBootSec=1w OnUnitInactiveSec=1w